# templates/beaver.service.erb
[Unit]
Description=Beaver log muncher for logstash

[Service]
ExecStart=/usr/local/bin/beaver -c ${BEAVER_CONFIG:-/etc/beaver.conf} -t redis $BEAVER_OPTIONS
Restart=on-failure
KillMode=process
