# templates/beaver.service.erb
<%
  if @provider == 'virtualenv'
    bin_path = "#{@venv}/bin/beaver"
  else
    bin_path = '/usr/local/bin/beaver'
  end
-%>
[Unit]
Description=Beaver log muncher for logstash

[Service]
ExecStart=<%= bin_path %> -c /etc/beaver/beaver.conf -t redis -l /var/log/beaver.log $BEAVER_OPTIONS
Restart=on-failure
KillMode=process

[Install]
WantedBy=multi-user.target
